# scripts/launchd/archived/com.noogram.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.noogram.dashboard</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/you/.pyenv/versions/3.10.15/bin/python</string>
    <string>/Users/you/dev/projects/noogram-dashboard/dashboard.py</string>
    <string>--port</string>
    <string>18790</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/you/dev/projects/noogram-dashboard</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/you/.pyenv/versions/3.10.15/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/noogram-dashboard.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/noogram-dashboard.stderr.log</string>
</dict>
</plist>
